# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing
ln  -s  /usr/bin/ntfs-3g   /sbin/mount.ntfs
/etc/init.d/xupnpd disable
exit 0
